# dist/xserver.conf
description "Play server"
start on vagrant-mounted

respawn
respawn limit unlimited

chdir /mnt/server/

setuid vagrant

env HOME=/home/vagrant/

script
	/usr/local/bin/npm start > /mnt/server/tmp/server.log
	sleep 5
end script
