# dist/xserver.conf
description "Play server"
start on vagrant-mounted

respawn
respawn limit unlimited

chdir /mnt/server/

setuid vagrant

env HOME=/home/vagrant
env NVM_DIR=/home/vagrant/.nvm

script
	/home/vagrant/.nvm/nvm-exec npm start 2>&1 > /mnt/server/tmp/server.log
	sleep 5
end script
