# dist/xserver.conf
description "Play server"
start on vagrant-mounted

respawn
respawn limit unlimited

setuid vagrant

script
	cd /mnt/server/
	export HOME=/home/vagrant
	export NVM_DIR=/home/vagrant/.nvm
	/home/vagrant/.nvm/nvm-exec npm start 2>&1 > /mnt/server/tmp/server.log
	sleep 5
end script
